# crontab/crontab
@hourly /usr/bin/docker system prune -f
